# usr/lib/systemd/system/go-aliddns.timer
[Unit]
Description="Go Aliddns Timer"

[Timer]
# 年、月、日、时、分、秒
OnCalendar=*-*-* 0/2:00:00
Persistent=true

[Install]
WantedBy=timers.target
